# configs/ble_gateway.service
[Unit]
After=bluetooth.service

[Service]
ExecStart=/usr/local/bin/ble_gateway.sh

[Install]
WantedBy=default.target
